# macbook/ansible/templates/com.prometheus.node_exporter.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!-- Managed by Ansible (macbook/ansible/playbooks/node-exporter.yaml). Do not edit by hand. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>{{ ne_label }}</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{ ne_bin }}</string>
        <string>--web.listen-address={{ ne_listen }}</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>{{ ansible_facts.user_dir }}/.local/var/log/node_exporter.log</string>
    <key>StandardErrorPath</key>
    <string>{{ ansible_facts.user_dir }}/.local/var/log/node_exporter.err.log</string>
</dict>
</plist>
